# recipes-ubnt/unifi/unifi/unifi.service.in
[Unit]
Description=Ubiquiti UniFi Controller Server
Wants=network-online.target time-sync.target
After=network-online.target time-sync.target

[Service]
User=@UNIFI_USER@
WorkingDirectory=@libdir@/unifi
RuntimeDirectory=unifi
LimitCORE=0
Type=notify
KillMode=control-group
NotifyAccess=main
ExecStart=@libdir@/unifi/bin/ace.sh start
SuccessExitStatus=143
RestartPreventExitStatus=2
TimeoutSec=240
Restart=on-failure

[Install]
WantedBy=multi-user.target
