# recipes-ubnt/unifi/unifi/unifi.service.in
[Unit]
Description=Ubiquiti UniFi Controller Server
Wants=network-online.target time-sync.target
After=network-online.target time-sync.target

[Service]
User=@UNIFI_USER@
WorkingDirectory=@libdir@/unifi
LimitCORE=0
EnvironmentFile=-@sysconfdir@/default/unifi.env
ExecStart=@bindir@/java $JAVA_ARGS -jar @libdir@/unifi/lib/ace.jar start
ExecStop=@bindir@/java $JAVA_ARGS -jar @libdir@/unifi/lib/ace.jar stop
SuccessExitStatus=143
TimeoutSec=240
Restart=on-failure

[Install]
WantedBy=multi-user.target
