# recipes-ubnt/unifi/unifi/unifi.service.in
[Unit]
Description=Ubiquiti UniFi Controller Server
Wants=network-online.target time-sync.target
After=network-online.target time-sync.target

[Service]
User=@UNIFI_USER@
WorkingDirectory=@libdir@/unifi
LimitCORE=0
ExecStart=@libdir@/unifi/bin/ace.sh start
ExecStop=@libdir@/unifi/bin/ace.sh stop
SuccessExitStatus=143
TimeoutSec=240
Restart=on-failure

[Install]
WantedBy=multi-user.target
